# scripts/parallel-sim-eod.service
[Unit]
Description=EOD Analysis — Paper vs Simulazione
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/htpc/backtrader/backtrader
EnvironmentFile=/home/htpc/backtrader/env/pa2
EnvironmentFile=/home/htpc/backtrader/env/psim
ExecStart=/home/htpc/backtrader/scripts/parallel-sim-eod.sh
StandardOutput=append:/home/htpc/backtrader/logs/parallel-sim-eod.log
StandardError=append:/home/htpc/backtrader/logs/parallel-sim-eod.log

[Install]
WantedBy=default.target
