# scripts/zmq-logger.service
[Unit]
Description=Alpaca ZMQ Flow Logger
After=network.target
Requires=zmq-proxy.service

[Service]
Type=simple
WorkingDirectory=/home/htpc/backtrader/bt-gateway/bt_alpaca_zmq

# Mostra nei log di systemd dove si trovano i file prodotti
ExecStartPre=/bin/bash -c 'echo "[ZMQ-LOGGER] scrivo i log in /home/htpc/backtrader/out/dump (file per simbolo e zmq_messages.log)"'

# Il logger non richiede le API ma eredita eventuali variabili dal file .zmq-env
ExecStart=/bin/bash -c 'source ../../env/zmq && exec ../.venv/bin/python ./zmq_logger.py'

StandardOutput=append:/home/htpc/backtrader/logs/zmq-logger.log
StandardError=append:/home/htpc/backtrader/logs/zmq-logger.log

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
